# dist/linux/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:app]
priority=1
startsecs = 0
command=/home/app
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
priority=2
startsecs = 0
stdout_logfile =/tmp/nginx.log
command=/bin/ash -c 'sleep 2 && /usr/sbin/nginx -g "daemon off;"'
